# docker/supervisord.conf
[supervisord]
nodaemon=true
user=www-data  ; Run Supervisor as www-data user

[program:php-fpm]
command=php-fpm
autostart=true
autorestart=true
stderr_logfile=/var/log/php-fpm.err.log
stdout_logfile=/var/log/php-fpm.out.log
user=www-data  ; Run PHP-FPM as www-data

[program:nginx]
command=nginx -g 'daemon off;'
autostart=true
autorestart=true
stderr_logfile=/var/log/nginx.err.log
stdout_logfile=/var/log/nginx.out.log
user=www-data  ; Run Nginx as www-data
